# openrc/init.d/acpid
#!/sbin/openrc-run
# Copyright (c) 2007-2008 Roy Marples <roy@marples.name>
# Released under the 2-clause BSD license.

command=/usr/sbin/acpid
pidfile=/var/run/acpid.pid
name="ACPI Daemon"

depend()
{
	after bootmisc
	use logger
	keyword -openvz -prefix -vserver -lxc
}

start_post()
{
	/bin/setterm -blank 15 -powersave powerdown -powerdown 60
}
