# openrc/init.d/atd
#!/sbin/openrc-run
# Copyright (c) 2007-2008 Roy Marples <roy@marples.name>
# Released under the 2-clause BSD license.

command=/usr/sbin/atd
command_args="-b 15 -l 1"
pidfile=/var/run/atd.pid

depend()
{
	need clock localmount
	after cron
}
