# openrc/init.d/avahi-dnsconfd
#!/sbin/openrc-run
# Copyright (c) 2007-2008 Roy Marples <roy@marples.name>
# Released under the 2-clause BSD license.

command=/usr/sbin/avahi-dnsconfd
pidfile=/var/run/avahi-dnsconfd.pid
extra_started_commands="reload"

depend()
{
  need avahid
  use logger dns net
  after bootmisc
}

start() {
  ebegin "Starting avahi-dnsconfd"
  $command -D
  eend $?
}

stop() {
  ebegin "Stopping avahi-dnsconfd"
  $command -k
  eend $?
}

reload() {
  ebegin "Reloading avahi-dnsconfd"
  $command -r
  eend $?
}
