# openrc/init.d/avahid
#!/sbin/openrc-run
# Copyright (c) 2007-2008 Roy Marples <roy@marples.name>
# Released under the 2-clause BSD license.

command=/usr/sbin/avahi-daemon
pidfile=/var/run/avahi-daemon/pid
extra_started_commands="reload"

depend()
{
	use net
	need dbus
	after bootmisc
}

start() {
  ebegin "Starting avahi-daemon"
  $command -D
  eend $?
}

stop() {
  ebegin "Stopping avahi-daemon"
  $command -k
  eend $?
}

reload() {
  ebegin "Reloading avahi-daemon"
  $command -r
  eend $?
}
