# openrc/init.d/consolefont.slack
#!/sbin/openrc-run
# Copyright (c) 2007-2008 Roy Marples <roy@marples.name>
# Released under the 2-clause BSD license.

description="Update /etc/mtab to match what the kernel knows about"

depend()
{
	need root
	keyword -prefix
}

start()
{
	ebegin "Setting up console font"
	setfont -v Lat2-Terminus16.psfu.gz
}
