# openrc/init.d/cron
#!/sbin/openrc-run
# Copyright (c) 2007-2008 Roy Marples <roy@marples.name>
# Released under the 2-clause BSD license.

command=/usr/sbin/crond
command_args="-l notice"
pidfile=/var/run/crond.pid

depend()
{
	need clock localmount
	after bootmisc
	provide cron
}

start_pre() {
    if ! ps axc | grep crond 1> /dev/null 2> /dev/null ; then
        rm -f ${pidfile}
    fi
}

start_post() {
    pidof $command > $pidfile
}

stop_post() {
    rm -f $pidfile
}
