# openrc/init.d/cupsd
#!/sbin/openrc-run
# Copyright (c) 2007-2008 Roy Marples <roy@marples.name>
# Released under the 2-clause BSD license.

command=/usr/sbin/cupsd
pidfile=/var/run/cups/cups.pid
name="CUPS Printing Daemon"

depend()
{
	need loopback logger
	after bootmisc
}

start()
{
	mkdir -p $(dirname $pidfile)
	rm -f $pidfile
	$command
	pidof cupsd > $pidfile
}
