# openrc/init.d/cupsd
#!/sbin/openrc-run
# Copyright (c) 2007-2008 Roy Marples <roy@marples.name>
# Released under the 2-clause BSD license.

command=/usr/sbin/cupsd
name="CUPS Printing Daemon"

depend()
{
	use net logger
	need dbus
	before nfs
	after logger
}

start_pre()
{
	mkdir -p /var/run/cups
}
