# openrc/init.d/dbus
#!/sbin/openrc-run
# Copyright (c) 2007-2008 Roy Marples <roy@marples.name>
# Released under the 2-clause BSD license.

command=/usr/bin/dbus-daemon
pidfile=/var/run/dbus/dbus.pid
command_args="${dbusd_args---system}"
name="Message Bus Daemon"

depend()
{
	need localmount syslogd
	after bootmisc syslogd
	use logger
}

start_pre()
{
	mkdir -p $(dirname $pidfile)
}
