# openrc/init.d/dbus
#!/sbin/openrc-run
# Copyright (c) 2007-2008 Roy Marples <roy@marples.name>
# Released under the 2-clause BSD license.

command=/usr/bin/dbus-daemon
pidfile=/var/run/dbus/dbus.pid
command_args="${dbusd_args---system}"
extra_started_commands="reload"
name="Message Bus Daemon"

depend()
{
	use logger
	after bootmisc
}

start_pre()
{
	mkdir -p $(dirname $pidfile)
}

reload() {
  ebegin "Reloading system message bus configuration..."
  if [ -e "$pidfile" ]; then
    pid=$(cat $pidfile)
    kill -HUP $pid
  else
    killall -HUP dbus-daemon
  fi
  eend $?
}

