# openrc/init.d/fontcache
#!/sbin/openrc-run
# Copyright (c) 2007-2008 Roy Marples <roy@marples.name>
# Released under the 2-clause BSD license.

name="X font cache"
command_background="yes"

depend() {
    need localmount ldconfig
    after bootmisc
}

start()
{
    ebegin "Updating X font cache"
    fc-cache -s -f 2>/dev/null &
    eend $?
}
