# openrc/init.d/glibschemas
#!/sbin/openrc-run
# Copyright (c) 2007-2008 Roy Marples <roy@marples.name>
# Released under the 2-clause BSD license.

name="Glib2 schemas compilation"
command_background="yes"

depend() {
    need localmount
    after bootmisc
}

start()
{
  ebegin "Updating glib2 schemas compilation"
  if [ -x /usr/bin/glib-compile-schemas ]; then
    /usr/bin/glib-compile-schemas /usr/share/glib-2.0/schemas >/dev/null 2>&1 &
  fi
  eend $?
}
