# openrc/init.d/inetd
#!/sbin/openrc-run
# Start/stop/restart inetd, the BSD Internet super-daemon.

description="inetd - the BSD Internet super-daemon"
command="/usr/sbin/inetd"
pidfile="/var/run/inetd.pid"

depend() {
    need net
    use logger
}

start() {
    start-stop-daemon --start --quiet --exec $command
}
