# openrc/init.d/ldconfig
#!/sbin/openrc-run

depend() {
    need localmount syslogd
    after bootmisc
}

start() {
    ebegin "Updating shared library links"
    /sbin/ldconfig &
    eend $?
}
