# openrc/init.d/loop
#!/sbin/openrc-run
# Copyright (c) 2007-2008 Roy Marples <roy@marples.name>
# Released under the 2-clause BSD license.

name="loop device"

depend()
{
    before modules
}

start()
{
	ebegin "Initializing loop device"
	/sbin/modprobe loop
	eend $?
}
