# openrc/init.d/mysqld
#!/sbin/openrc-run

name="mysqld"
description="MySQL Server"
pidfile="/run/mysql/mysql.pid"
command="/usr/bin/mysqld_safe"
command_background="yes"
command_args="--user=mysql --datadir=/var/lib/mysql --pid-file=$pidfile $SKIP $INNODB"

depend() {
    use net.lo
    need localmount loopback
}

start_pre() {
    if [ -r $pidfile ]; then
      if ! ps axc | grep mysqld 1> /dev/null 2> /dev/null ; then
        rm -f ${pidfile}
      fi
    else
      mkdir -p $(dirname $pidfile)
      chown -R $user:$group $(dirname $pidfile)
    fi
}

stop() {
  # If there is no PID file, ignore this request...
  if [ -r /var/run/mysql/mysql.pid ]; then
    killall mysqld
    # Wait at least one minute for it to exit, as we don't know how big the DB is...
    for second in 0 1 2 3 4 5 6 7 8 9 0 1 2 3 4 5 6 7 8 9 0 1 2 3 4 5 6 7 8 9 \
      0 1 2 3 4 5 6 7 8 9 0 1 2 3 4 5 6 7 8 9 0 1 2 3 4 5 6 7 8 9 60 ; do
      if [ ! -r /var/run/mysql/mysql.pid ]; then
        break;
      fi
      sleep 1
    done
    if [ "$second" = "60" ]; then
      echo "WARNING:  Gave up waiting for mysqld to exit!"
      sleep 15
    fi
  fi
}
