# openrc/init.d/NetworkManager
#!/sbin/openrc-run
# Copyright (c) 2008 Saleem Abdulrasool <compnerd@compnerd.org>
# Distributed under the terms of the GNU General Purpose License v2
# $Header: /var/cvsroot/gentoo-x86/net-misc/networkmanager/files/init.d.NetworkManager,v 1.1 2013/01/28 07:05:05 tetromino Exp $

description="NetworkManager daemon. The service is marked as started only \
when a network connection is established."

command=/usr/sbin/NetworkManager
pidfile=/var/run/NetworkManager/NetworkManager.pid
name="NetworkManager Daemon"

depend() {
	provide net
	need localmount dbus
	use logger
}

start()
{
  # Just in case the pidfile is still there, we may need to nuke it.
  if [ -e "$pidfile" ]; then
    rm -f $pidfile
  fi
  $command
}

stop()
{
  # Shut down any DHCP connections, otherwise the processes will be orphaned
  # and the connections will not come up when NetworkManager restarts.
  #if ps ax | grep /sbin/dhcpcd | grep -q libexec/nm-dhcp ; then
    #ps ax | grep /sbin/dhcpcd | grep libexec/nm-dhcp | while read line ; do
      #kill -HUP $(echo $line | cut -b 1-5)
    #done
  #fi
  #if ps ax | grep /sbin/dhclient | grep -q /var/lib/NetworkManager ; then
    #ps ax | grep /sbin/dhclient | grep /var/lib/NetworkManager | while read line ; do
      #kill -HUP $(echo $line | cut -b 1-5)
    #done
  #fi
  kill `cat $pidfile 2>/dev/null`
  sleep 3
  rm -f $pidfile &>/dev/null
}

restart()
{
  nm_stop
  nm_start
}
