# openrc/init.d/php-fpm
#!/sbin/openrc-run

name="php-fpm"
description="PHP FastCGI Process Manager"
extra_started_commands="reload"
pidfile="/var/run/php-fpm.pid"

depend() {
        need net
}

start() {
        ebegin "Starting PHP FastCGI server"
        start-stop-daemon --start --pidfile=$pidfile \
               --exec /usr/sbin/php-fpm
        eend $?
}

stop() {
        ebegin "Stopping PHP FastCGI server"
        start-stop-daemon --stop --pidfile=$pidfile
        eend $?
}

reload() {
        ebegin "Restart gracefully PHP FastCGI server"
        kill -USR2 `cat $pidfile` &>/dev/null
        eend $?
}
