# openrc/init.d/php-fpm
#!/sbin/openrc-run

name="php-fpm"
description="PHP FastCGI Process Manager"
command="/usr/sbin/php-fpm"
command_args="-y /etc/php-fpm.conf"
extra_started_commands="reload"
pidfile="/var/run/php-fpm.pid"

depend() {
        need net
}

reload() {
        ebegin "Restart gracefully PHP FastCGI server"
        kill -USR2 `cat $pidfile` &>/dev/null
        eend $?
}
