# openrc/init.d/rsyslog
#!/sbin/openrc-run

command=/usr/sbin/rsyslogd
name="enhanced syslogd"

depend()
{
    provide rsyslogd syslog
    need $remote_fs $time
}
