# openrc/init.d/snmpd
#!/sbin/runscript
# Copyright 1999-2012 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-analyzer/net-snmp/files/snmpd.init.2,v 1.3 2012/10/22 02:57:05 flameeyes Exp $

SNMPD_PIDFILE="${SNMPD_PIDFILE:-/var/run/snmpd}"
SNMPD_CONFIG="${SNMPD_CONFIG:-/etc/snmp/snmpd.conf}"
extra_started_commands="reload"

command="/usr/sbin/snmpd"
command_args="-A -p ${SNMPD_PIDFILE} -a ${SNMPD_FLAGS} -c ${SNMPD_CONFIG}"
pidfile="${SNMPD_PIDFILE}"

depend() {
	use logger
}

checkconfig() {
	if [ ! -e ${SNMPD_CONFIG} ] ; then
		eerror "${SVCNAME} requires an ${SNMPD_CONFIG} configuration file"
		return 1
	fi
}

start_pre() {
	checkconfig || return 1
}

stop() {
    ebegin "Stopping snmpd"
    kill -9 `cat $pidfile` 2>/dev/null
    sleep 1
    rm $pidfile
    eend $? "Failed to stop snmpd"
}

reload() {
	checkconfig || return 1

	ebegin "Reloading ${SVCNAME} configuration"
	kill -HUP $(cat ${SNMPD_PIDFILE}) 2>&1 > /dev/null
	eend $?
}
