# openrc/init.d/xdm
#!/sbin/openrc-run
# Copyright 2015 walecha.web.id
# Distributed under the terms of the GNU General Public License, v2 or later

name="Starting X11 Display Manager"
DM_USE=${DM_USE:-sddm}

depend() {
    need localmount
    use dbus consolekit elogind
    after bootmisc
}

start() {
    ebegin "Starting display manager"
    case $DM_USE in
	sddm)
	    mkdir -p /run/sddm
	    chown -R sddm:sddm /run/sddm
	    start-stop-daemon --background -m -p /run/sddm/sddm.pid --start --exec /usr/bin/sddm
	    eend $?
	    ;;
	gdm)
	    mkdir -p /run/gdm
	    chown -R gdm:gdm /run/gdm
	    start-stop-daemon --background --start --exec /usr/sbin/gdm
	    eend $?
	    ;;
	lxdm)
	    mkdir -p /run/lxdm
	    chown -R nobody:nobody /run/lxdm
	    start-stop-daemon --background --start --exec /usr/bin/lxdm
	    eend $?
	    ;;
	xdm)
	    mkdir -p /run/xdm
	    chown -R nobody:nobody /run/xdm
	    start-stop-daemon --background --start --exec /usr/bin/xdm
	    eend $?
	    ;;
    esac
}

stop() {
    ebegin "Stopping display manager"
    case $DM_USE in
	sddm)
	    mkdir -p /run/sddm
	    chown -R sddm:sddm /run/sddm
	    start-stop-daemon --stop --exec /usr/bin/sddm && sleep 3 && killall sddm >/dev/null
	    eend $?
	    ;;
	gdm)
	    mkdir -p /run/gdm
	    chown -R gdm:gdm /run/gdm
	    start-stop-daemon --stop --exec /usr/sbin/gdm && sleep 3 && killall gdm >/dev/null
	    eend $?
	    ;;
	lxdm)
	    mkdir -p /run/lxdm
	    chown -R nobody:nobody /run/lxdm
	    start-stop-daemon --stop --exec /usr/bin/lxdm && sleep 3 && killall lxdm >/dev/null
	    eend $?
	    ;;
	xdm)
	    mkdir -p /run/xdm
	    chown -R nobody:nobody /run/xdm
	    start-stop-daemon --stop --exec /usr/bin/xdm && sleep 3 && killall xdm >/dev/null
	    eend $?
	    ;;
    esac
}
