# openrc/init.d.old/cgroup-mount
#!/sbin/openrc-run

depend()
{
	# we depend on udev-mount explicitly, not dev-mount generic as we don't
	# want mdev as a dev-mount provider to come in.
	need sysfs udev-mount
	before checkfs fsck

	# udev does not work inside vservers
	keyword -vserver -lxc
}

start() {
    # Mount Control Groups filesystem interface:
    if grep -wq cgroup /proc/filesystems ; then
        if [ -d /sys/fs/cgroup ]; then
            # See linux-*/Documentation/cgroups/cgroups.txt (section 1.6)
            # Check if we have some tools to autodetect the available cgroup controllers
            if [ -x /bin/cut -a -x /bin/tail ]; then
            # Mount a tmpfs as the cgroup filesystem root
            mount -t tmpfs -o mode=0755,size=8M cgroup_root /sys/fs/cgroup
            # Autodetect available controllers and mount them in subfolders
            controllers="$(/bin/cut -f 1 /proc/cgroups | /bin/tail -n +2)"
            for i in $controllers; do
                mkdir /sys/fs/cgroup/$i
                mount -t cgroup -o $i $i /sys/fs/cgroup/$i
            done
            unset i controllers
            else
            # We can't use autodetection so fall back mounting them all together
            mount -t cgroup cgroup /sys/fs/cgroup
            fi
        else
            mkdir -p /dev/cgroup
            mount -t cgroup cgroup /dev/cgroup
        fi
    fi
}
