# openrc/init.d.old/udev-settle
#!/sbin/openrc-run
# This is a script to initialize udev, which populates the /dev
# directory with device nodes, scans for devices, loads the
# appropriate kernel modules, and configures the devices.

depend()
{
	use root
	before localmount

	# udev does not work inside vservers
	keyword -vserver -lxc
}

#PATH="/sbin:/bin"

. /etc/udev/udev.conf

# remove trailing slash from udev_root
UDEV_ROOT=$(echo "${udev_root}" |sed 's/\/*$//')

start() {
    # If the persistent rules files do not exist, trigger an add event:
    if [ ! -r /etc/udev/rules.d/70-persistent-net.rules -o ! -r /etc/udev/rules.d/70-persistent-cd.rules ]; then
	# Test that we can actually write to the directory first:
	if touch /etc/udev/rules.d/testfile 2> /dev/null ; then
	    rm -f /etc/udev/rules.d/testfile
	    # This should add persistent net/cd rules:
	    /sbin/udevadm trigger --type=devices --action=add
	    sleep 3
	    # Create the files if they don't exist at this point.
	    # If a machine does not have a network device or an optical
	    # device, we don't want to waste time trying to generate
	    # rules at every boot.
	    # To force another attempt, delete the file(s).
	    touch /etc/udev/rules.d/70-persistent-net.rules
	    touch /etc/udev/rules.d/70-persistent-cd.rules
	fi
	# Since udevd is running, most of the time we only need change events:
	/sbin/udevadm trigger --type=subsystems --action=change
	/sbin/udevadm trigger --type=devices --action=change
    fi
    /sbin/udevadm settle --timeout=120
}
